# genesi-arch/packages/genesi-updater/genesi-updater.service
[Unit]
Description=Genesi OS Update Checker
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/genesi-updater
User=nobody
Group=nobody

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/genesi-updater /var/log /tmp

[Install]
WantedBy=multi-user.target
