# genesi-arch/packages/genesi-updater/genesi-updater.timer
[Unit]
Description=Genesi OS Update Checker Timer
Requires=genesi-updater.service

[Timer]
# Check for updates every hour
OnBootSec=5min
OnUnitActiveSec=1h

# Randomize by 10 minutes to avoid server load spikes
RandomizedDelaySec=10min

# Persistent timer (runs missed checks on boot)
Persistent=true

[Install]
WantedBy=timers.target
